# deploy/systemd/shkila-payment-sync.service
[Unit]
Description=Shkila YooKassa pending payments reconciliation
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/var/www/shkila/back/school
ExecStart=/var/www/shkila/venv/bin/python manage.py sync_pending_payments --max-age-hours 72
User=root
Group=root

